# config/includes.chroot_after_packages/home/user/.config/autostart/dosbox.desktop
[Desktop Entry]
Type=Application
Name=Dosbox
Exec=/usr/bin/dosbox
